# systemd_file_kernel_log_backup.timer
[Unit]
Description=Hourly Kernel Log Backup Timer

[Timer]
OnCalendar=*:0
Persistent=true

[Install]
WantedBy=timers.target
